# pkg/sagan.rc
#!/bin/ksh
#
# $OpenBSD: sagan.rc,v 1.0 2020/02/05 01:03:23 litew Exp $

daemon="/usr/local/bin/sagan"
daemon_flags="-D"

rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
